# dotfiles/.config/autostart/hyperhdr.desktop
[Desktop Entry]
Type=Application
Name=HyperHDR
Exec=hyperhdr
