# dotfiles/.config/autostart/instant-eyedropper.desktop
[Desktop Entry]
Type=Application
Name=Instant Eyedropper Reborn
Exec=ie-r
