# dotfiles/.config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Type=Application
Name=Pipewire Pulse
Exec=pipewire-pulse
